# dbgpt.service
# DB-GPT Service Unit File
# 服务管理命令：
#   启动服务：sudo systemctl start dbgpt
#   停止服务：sudo systemctl stop dbgpt
#   重启服务：sudo systemctl restart dbgpt
#   查看状态：sudo systemctl status dbgpt
#   查看日志：sudo journalctl -u dbgpt -f
#   启用开机启动：sudo systemctl enable dbgpt
#   禁用开机启动：sudo systemctl disable dbgpt

[Unit]
Description=DB-GPT Service - 数据库智能对话服务
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=adminstrator
WorkingDirectory=/home/adminstrator/work/Finai
# 设置运行时环境变量
Environment="PATH=/home/adminstrator/.local/bin:/home/adminstrator/work/Finai/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/adminstrator"
# 服务启动命令
ExecStart=/home/adminstrator/.local/bin/uv run python packages/dbgpt-app/src/dbgpt_app/dbgpt_server.py --config configs/dbgpt-proxy-siliconflow.toml
# 故障重启策略
Restart=on-failure
RestartSec=10
# 日志输出配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dbgpt

# 自动重启策略：300秒内如果重启超过5次则停止自动重启
StartLimitInterval=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
